# packer/include/runsvdir.conf
start on runlevel [2345]
stop on runlevel [016]
respawn
kill signal HUP
exec /usr/sbin/runsvdir-start